# scripts/ai.openclaw.taskai-whatsapp-dispatch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.openclaw.taskai-whatsapp-dispatch</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/bobbylian/.local/bin/taskai-whatsapp-dispatch-run.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/bobbylian/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/tmp/taskai-whatsapp-dispatch.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/taskai-whatsapp-dispatch.err</string>
</dict>
</plist>
